# ex04/10-module_keyboard.rules
ACTION="add", SUBSYSTEMS=="usb", ATTR{bInterfaceProtocol}=="01", RUN+="/sbin/insmod /home/mhouppin/little-penguin/ex04/main.ko"